# ansible/roles/lan/files/route53-hostname.timer
[Unit]
Description=Update Route53 record timer

[Timer]
OnBootSec=1m
OnUnitInactiveSec=2m

[Install]
WantedBy=timers.target
